# deploy/systemd/durak-bot.service
[Unit]
Description=Durak Telegram Bot Worker
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/var/www/app.games-telegram.online/current/backend
EnvironmentFile=/var/www/app.games-telegram.online/shared/config/bot.env
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=3
StandardOutput=append:/var/www/app.games-telegram.online/shared/logs/bot.log
StandardError=append:/var/www/app.games-telegram.online/shared/logs/bot.log

[Install]
WantedBy=multi-user.target

