# deploy/systemd/durak-backend.service
[Unit]
Description=Durak Telegram Backend
After=network.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/var/www/app.games-telegram.online/current/backend
EnvironmentFile=/var/www/app.games-telegram.online/shared/config/backend.env
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=3
StandardOutput=append:/var/www/app.games-telegram.online/shared/logs/backend.log
StandardError=append:/var/www/app.games-telegram.online/shared/logs/backend.log

[Install]
WantedBy=multi-user.target

